# pi/heartwood-bunker.service
[Unit]
Description=Heartwood NIP-46 bunker
After=network-online.target heartwood.service
Wants=network-online.target
Requires=heartwood.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=heartwood
Group=heartwood
WorkingDirectory=/opt/heartwood/bunker
ExecStart=/usr/bin/node index.mjs
Restart=always
RestartSec=15
Environment=NODE_ENV=production

# Security hardening
ProtectSystem=strict
ProtectHome=true
# Note: MemoryDenyWriteExecute is intentionally omitted — Node.js (V8 JIT)
# requires write-execute memory pages. The Rust binary unit includes it.
ReadWritePaths=/var/lib/heartwood
ReadOnlyPaths=/opt/heartwood
NoNewPrivileges=true
PrivateTmp=true

# Capability and syscall restrictions
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
